# deploy/clewdr.service
[Unit]
Description=clewdr-hub - Claude shared gateway
After=network.target

[Service]
Type=simple
User=clewdr
Group=clewdr
WorkingDirectory=/opt/clewdr
ExecStart=/opt/clewdr/clewdr --log-dir /opt/clewdr/log
Restart=on-failure
RestartSec=5

Environment=CLEWDR_IP=0.0.0.0
Environment=CLEWDR_PORT=8484

StandardOutput=journal
StandardError=journal
SyslogIdentifier=clewdr

[Install]
WantedBy=multi-user.target
